# deploy/broombuster.service
# systemd unit: BroomBuster API (no-login) on an always-on Linux host.
# Installed by deploy/install-service.sh (fills in user + repo path).
# Tailscale HTTPS is set up once, separately: `tailscale serve --bg 8000`.
[Unit]
Description=BroomBuster (no-login street-sweeping map)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
User=__USER__
WorkingDirectory=__REPO__
Environment=DEV_MODE=true
ExecStart=__REPO__/.venv/bin/python -m uvicorn broombuster.api.app:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
